# deploy/com.hanserautisch.investment-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.hanserautisch.investment-backup</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/nilskrueger/LUNA-Backups/backup-from-nas.sh</string>
  </array>
  <key>WorkingDirectory</key><string>/Users/nilskrueger/LUNA-Backups</string>
  <key>EnvironmentVariables</key>
  <dict><key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string></dict>
  <key>StartCalendarInterval</key>
  <dict><key>Hour</key><integer>3</integer><key>Minute</key><integer>0</integer></dict>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/Users/nilskrueger/Library/Logs/investment-backup.log</string>
  <key>StandardErrorPath</key><string>/Users/nilskrueger/Library/Logs/investment-backup.log</string>
</dict>
</plist>
